# deploy/systemd/geegoo-agent-pre-market.timer
[Unit]
Description=GeeGoo Agent pre-market timer (Mon-Fri 08:00 Asia/Shanghai)
[Timer]
OnCalendar=Mon..Fri *-*-* 08:00:00
Timezone=Asia/Shanghai
Persistent=true
Unit=geegoo-agent-pre-market.service

[Install]
WantedBy=timers.target
